# config/supervisord.d/dnsmasq.ini
[program:dnsmasq]
command = dnsmasq
    --no-daemon
    --port=19321
    --no-resolv
    --address=/hostmachine.dev/%(ENV_HOSTMACHINE_IP)s
    --address=/.nginxdemo/%(ENV_DNS_IP)s
stdout_events_enabled = true
stderr_events_enabled = true
autostart=true
autorestart=true
